# docker/supervisor/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=2
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout-nginx.log
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr-nginx.log
stderr_logfile_maxbytes=0
